# root/etc/periodic/15min/nextcloud
#!/bin/sh
sudo -u nextcloud php -f /opt/nextcloud/cron.php > /dev/stdout 2>&1
