# root/etc/periodic/15min/nextcloud
#!/bin/sh
sudo -u nextcloud php7 -f /opt/nextcloud/cron.php > /dev/null 2>&1
